# deployment/rich-list.service
[Unit]
Description=Rich-List Node.js Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=jim
WorkingDirectory=/opt/rich-list
EnvironmentFile=/opt/rich-list/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rich-list

[Install]
WantedBy=multi-user.target